# rf-backup@.service
[Unit]
Description=RF backup service

[Service]
Type=oneshot
ExecStart=/usr/share/rf-backup/do-backup.sh %i
